# deploy/systemd/financial-manager-webui.service
# 安装前：将下文所有 /opt/FinancialManager 改为你的项目根目录；
# 将 User/Group 改为与 fm-api.service 一致。
# 同机调 API 时 .env 中建议：FM_API_BASE=http://127.0.0.1:8123/api/v1（端口与 fm-api 一致）
#
# sudo cp fm-webui.service /etc/systemd/system/
# sudo systemctl daemon-reload
# sudo systemctl enable --now fm-webui.service

[Unit]
Description=Financial Manager Streamlit UI
After=network-online.target fm-api.service
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/FinancialManager
EnvironmentFile=/opt/FinancialManager/.env
Environment=PYTHONUNBUFFERED=1
Environment=STREAMLIT_BROWSER_GATHER_USAGE_STATS=false
ExecStart=/opt/FinancialManager/.venv/bin/streamlit run finance_webui.py --server.address=127.0.0.1 --server.port=8501 --server.headless=true
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
